# systemd/mod-host.service
[Unit]
Description=MOD-Host LV2 Plugin Host (JACK2 Native)
After=jackd.service
Requires=jackd.service
Documentation=man:mod-host(1)

[Service]
Type=simple
User=pilal
Group=pilal
SupplementaryGroups=audio

# Variables d'environnement
Environment="LV2_PATH=/var/modep/lv2:/usr/local/lib/lv2:/usr/lib/lv2"
Environment="MOD_LOG=1"

# Démarrage mod-host (JACK natif, sans pw-jack)
ExecStart=/usr/local/bin/mod-host -v -n -p 5555 -f 5556

# Redémarrage automatique
Restart=always
RestartSec=3

# Logs
StandardOutput=journal
StandardError=journal

# Working directory
WorkingDirectory=/var/modep

[Install]
WantedBy=multi-user.target
