# systemd/mod-host.service
[Unit]
Description=MOD-Host LV2 Plugin Host (via PipeWire-JACK)
After=sound.target
Documentation=man:mod-host(1)

[Service]
Type=simple
User=pilal
Group=pilal
SupplementaryGroups=audio

# Variables d'environnement pour PipeWire
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="LV2_PATH=/var/modep/lv2:/usr/local/lib/lv2:/usr/lib/lv2"
Environment="MOD_LOG=1"

# Démarrage mod-host via pw-jack (PipeWire JACK wrapper)
ExecStart=/usr/bin/pw-jack /usr/local/bin/mod-host -v -n -p 5555 -f 5556

# Redémarrage automatique
Restart=always
RestartSec=3

# Logs
StandardOutput=journal
StandardError=journal

# Working directory
WorkingDirectory=/var/modep

[Install]
WantedBy=multi-user.target
