# systemd/mod-ui.service
[Unit]
Description=MOD-UI Web Interface (JACK2 Native)
After=mod-host.service jackd.service
Requires=mod-host.service
Documentation=https://github.com/moddevices/mod-ui

[Service]
Type=simple
User=pilal
Group=pilal
SupplementaryGroups=audio

# Capability pour bind port 80 sans root
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Variables d'environnement Python
Environment="PYTHONPATH=/usr/local/lib/python3.10/site-packages"

# Variables d'environnement MOD-UI
Environment="MOD_DATA_DIR=/var/modep"
Environment="MOD_HTML_DIR=/usr/local/share/mod/html"
Environment="MOD_DEVICE_WEBSERVER_PORT=80"
Environment="LV2_PATH=/var/modep/lv2:/usr/local/lib/lv2:/usr/lib/lv2"
Environment="MOD_LOG=1"
Environment="MOD_USER_FILES_DIR=/var/modep/user-files"
Environment="MOD_PEDALBOARDS_DIR=/var/modep/pedalboards"
Environment="MOD_KEYS_PATH=/var/modep/data/keys"
Environment="MOD_BANKS_JSON_FILE=/var/modep/data/banks.json"
Environment="MOD_DEFAULT_PEDALBOARD=/usr/local/share/mod/default.pedalboard"

# Attente que mod-host soit prêt
ExecStartPre=/bin/sleep 5

# Démarrage mod-ui (JACK natif, sans pw-jack)
ExecStart=/usr/local/bin/mod-ui

# Redémarrage automatique
Restart=always
RestartSec=3

# Logs
StandardOutput=journal
StandardError=journal

# Working directory
WorkingDirectory=/var/modep

[Install]
WantedBy=multi-user.target
